# service/mail_send.timer
[Unit]
Description=mail sender
[Timer]
OnCalendar=*-*-* 18:45:00
[Install]
WantedBy=multi-user.target